# cron/crontab
#!/bin/bash
cd /usr/lib/zabbix/externalscripts

#CloudWatchLogs
for i in `aws logs describe-log-groups | jq -r '.logGroups[].logGroupName'`; do MSG=`/usr/lib/zabbix/externalscripts/cw_log_filter.rb "$i" ERR 300`; zabbix_sender -z localhost -p 10051 -s "AWS-TEST" -k "LogGroup[$i]" -o "$MSG";  done

#SQS Message Count
for i in `aws sqs list-queues | jq -r ".QueueUrls[]"`;do ApproximateNumberOfMessages=`aws sqs get-queue-attributes --queue-url $i --attribute-names ApproximateNumberOfMessages | jq .Attributes.ApproximateNumberOfMessages`; zabbix_sender -z localhost -p 10051 -s "AWS-TEST" -k "MSG_COUNT[$i]" -o "$ApproximateNumberOfMessages"; done


